# services/model-broker/deploy/openphone-model-broker.service
[Unit]
Description=OpenPhone Model Broker
Documentation=https://github.com/openphone-os/OpenPhone
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=openphone-broker
Group=openphone-broker
WorkingDirectory=/opt/openphone
EnvironmentFile=/etc/openphone/model-broker.env
ExecStart=/usr/bin/python3 /opt/openphone/services/model-broker/openphone_model_broker.py
Restart=on-failure
RestartSec=5s

NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ProtectClock=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictNamespaces=true
LockPersonality=true
MemoryDenyWriteExecute=true
SystemCallArchitectures=native
ReadWritePaths=/var/log/openphone
ReadOnlyPaths=/opt/openphone /etc/openphone
UMask=0077

[Install]
WantedBy=multi-user.target
